# config/etc-darwin/profile.d/nix-darwin.sh
export DYLD_FRAMEWORK_PATH=/System/Library/Frameworks
