# config/etc/profile.d/nix-johbo.sh

NIX_LINK=$HOME/.nix-profile

# use nix configuration out of my nix profile
export NIX_CONF_DIR=$NIX_LINK/etc/nix
