# config/etc-linux/supervisor.d/emacs-server.conf

[program:emacs-server]
command=daemon-fg --pid-file=~/var/run/emacs-server.pid -- emacs --daemon
autostart=true
