# config/etc/supervisor.d/emacs-server.conf

[program:emacs-server]
command=pidproxy ~/var/run/emacs-server.pid ~/.nix-profile/Applications/Emacs.app/Contents/MacOS/Emacs --daemon
autostart=false
