# config/etc/supervisor.d/emacs-server.conf

[program:emacs-server]
command=emacs-daemon.sh
autostart=false
